# configs/batteria/px13-battery-charge-limit.service
[Unit]
Description=ASUS PX13 battery charge-limit workaround
Documentation=https://bugs.kde.org/show_bug.cgi?id=450551

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'for _ in 1 2 3 4 5 6 7 8 9 10 11 12 13 14 15; do if [ -w /sys/class/power_supply/BAT0/charge_control_end_threshold ]; then printf "85\n" > /sys/class/power_supply/BAT0/charge_control_end_threshold; exit 0; fi; sleep 1; done; exit 1'
TimeoutStartSec=20
