# configs/audio-tas2783/px13-audio-resume.service
[Unit]
Description=ASUS PX13 TAS2783 audio resume workaround
After=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/px13-audio-resume
TimeoutStartSec=60

[Install]
WantedBy=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target
